# koala.service
[Unit]
Description=Koala
After=network.target

[Service]
WorkingDirectory=/opt/koala
ExecStart=/opt/koala/koala

[Install]
WantedBy=default.target
