# koala.service
[Unit]
Description=Koala
After=network.target

[Service]
WorkingDirectory=/koala
ExecStart=/koala/koala

[Install]
WantedBy=multi-user.target
